# drivers/scsi/qla2xxx/extras/99-qla2xxx.rules
#
# RULES FILE: located in /etc/udev/rules.d/
#

# qla2xxx driver
SUBSYSTEM=="pci", ENV{DRIVER}=="qla2xxx", ACTION=="change", RUN+="qla2xxx_udev.sh"
